# anacron/etc/anacrontab
# anacrontab: configuration file for anacron

# See anacron(8) and anacrontab(5) for details.

SHELL=/bin/sh
PATH=$HOME/bin:/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=snnw

1         5       cron.daily       nice run-parts $HOME/.anacron/cron.daily
7         10      cron.weekly      nice run-parts $HOME/.anacron/cron.weekly
@monthly  15      cron.monthly     nice run-parts $HOME/.anacron/cron.monthly
